# weather-station/tools/pi-weather-update.service
[Unit]
Description=Pi Weather Station auto-update (git pull from GitHub)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=flightrack
ExecStart=/home/flightrack/Dad-Flight-Tracker/weather-station/tools/auto-update.sh
